# integrations/wazuh-ioc-enrichment/misp-enricher.service
[Unit]
Description=Wazuh MISP IOC Webhook Enricher
After=network.target

[Service]
Type=simple
User=wazuh
Group=wazuh
NoNewPrivileges=yes
ProtectSystem=strict
PrivateTmp=yes
ReadWritePaths=/var/ossec/integrations/ioc
ReadWritePaths=/var/wazuh-ioc
WorkingDirectory=/var/ossec/integrations
ExecStart=/usr/bin/env python3 /var/ossec/integrations/misp_enricher.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
